# roles/nextcloud_client/templates/nextcloud-metered-monitor.service.j2
# {{ ansible_managed }}
[Unit]
Description=Nextcloud Metered Network Monitor
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart={{ nextcloud_client_script_dir }}/{{ nextcloud_client_script_name }}
Restart=always
RestartSec=10
Environment=HOME=%h
Environment=XDG_RUNTIME_DIR=/run/user/%i

# Security hardening
PrivateTmp=true
ProtectHome=false
ProtectSystem=strict
NoNewPrivileges=true
MemoryDenyWriteExecute=false
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
SystemCallArchitectures=native

# Allow network access and basic system calls
SystemCallFilter=@system-service
SystemCallFilter=~@debug @mount @cpu-emulation @obsolete @privileged @reboot @swap

[Install]
WantedBy=default.target
